# service-files/servo-zeroconf.service
[Unit]
Description=zeroconf

[Service]
ExecStart=/bin/su pi -c "/home/pi/vent/advertise.py vent-servo _http._tcp 80"
